# trunk/scripts/xxxopr.crontab
15 9 * * * cd /reg/g/pcds/dist/pds/tools/xtcclean; ./xtcclean2.py -t xpp -a --force | mail -s 'XPP xtcclean' weaver
